# debian/zabbix-agent-extra-smartctl.postrm
#!/bin/sh

set -e

case $1 in 
    remove|purge)
	invoke-rc.d zabbix-agent restart
esac

#DEBHELPER#
